# openrc/arbor
#!/sbin/openrc-run

description="Arbor web UI"
command="/usr/bin/arbor"
command_user="arbor"
command_background=true
pidfile="/run/arbor/web.pid"
output_log="/var/log/arbor/web.log"
error_log="/var/log/arbor/web.log"

env_file="/etc/arbor/arbor.env"

depend() {
    need arbor-daemon
    after net
}

start_pre() {
    checkpath --directory --owner root:arbor --mode 0750 /run/arbor
    checkpath --directory --owner arbor:arbor --mode 0755 /var/log/arbor
}
