# openrc/arbor-daemon
#!/sbin/openrc-run

description="Arbor privilege daemon"
command="/usr/bin/arbor-daemon"
command_background=true
supervisor="supervise-daemon"
respawn_delay=5
# Must be >= respawn_delay * respawn_max, otherwise supervise-daemon warns.
respawn_period=60
respawn_max=10
pidfile="/run/arbor/daemon.pid"
output_log="/var/log/arbor/daemon.log"
error_log="/var/log/arbor/daemon.log"
env_file="/etc/arbor/arbor.env"

depend() {
    need localmount
    after logger
}

start_pre() {
    checkpath --directory --owner root:arbor --mode 0750 /run/arbor
    checkpath --directory --owner arbor:arbor --mode 0750 /var/log/arbor
    touch /var/log/arbor/daemon.log
    chown root:arbor /var/log/arbor/daemon.log
    chmod 0640 /var/log/arbor/daemon.log
}
