# openrc/arbor
#!/sbin/openrc-run

description="Arbor web UI"
command="/usr/bin/arbor"
command_user="arbor"
command_background=true
supervisor="supervise-daemon"
respawn_delay=5
respawn_period=30
respawn_max=10
pidfile="/run/arbor/web.pid"
output_log="/var/log/arbor/web.log"
error_log="/var/log/arbor/web.log"

env_file="/etc/arbor/arbor.env"

depend() {
    need arbor-daemon
    after net
}

start_pre() {
    checkpath --directory --owner root:arbor --mode 0750 /run/arbor
    checkpath --directory --owner arbor:arbor --mode 0750 /var/log/arbor
    touch /var/log/arbor/web.log
    chown arbor:arbor /var/log/arbor/web.log
    chmod 0640 /var/log/arbor/web.log
}
