# openrc/arbor
#!/sbin/openrc-run

description="Arbor web UI"

# supervise-daemon's command_user handles the privilege drop to the arbor
# system user in all cases. If setpriv is also present (USE=setpriv in
# sys-apps/util-linux) we layer --no-new-privs on top so the process
# cannot regain capabilities via setuid binaries.
command_user="arbor:arbor"
_SETPRIV=/usr/bin/setpriv
if [ -x "$_SETPRIV" ]; then
    command="$_SETPRIV"
    command_args="--no-new-privs -- /usr/bin/arbor"
else
    command="/usr/bin/arbor"
    command_args=""
fi

command_background=true
supervisor="supervise-daemon"
respawn_delay=5
# Must be >= respawn_delay * respawn_max, otherwise supervise-daemon warns.
respawn_period=60
respawn_max=10
pidfile="/run/arbor/web.pid"
output_log="/var/log/arbor/web.log"
error_log="/var/log/arbor/web.log"

env_file="/etc/arbor/arbor.env"

depend() {
    need arbor-daemon
    after net
}

start_pre() {
    checkpath --directory --owner root:arbor --mode 0750 /run/arbor
    checkpath --directory --owner arbor:arbor --mode 0750 /var/log/arbor
    touch /var/log/arbor/web.log
    chown arbor:arbor /var/log/arbor/web.log
    chmod 0640 /var/log/arbor/web.log
    if [ ! -x /usr/bin/setpriv ]; then
        ewarn "setpriv not found — no_new_privs not applied."
        ewarn "For full hardening: USE=setpriv emerge sys-apps/util-linux"
    fi
}
