# openrc/arbor
#!/sbin/openrc-run

description="Arbor web UI"

# Hardening: drop to the unprivileged arbor user and lock no_new_privs so
# the web process cannot regain capabilities by exec'ing a setuid binary.
command="/usr/bin/setpriv"
command_args="--reuid=arbor --regid=arbor --init-groups --no-new-privs -- /usr/bin/arbor"

command_background=true
supervisor="supervise-daemon"
respawn_delay=5
# Must be >= respawn_delay * respawn_max, otherwise supervise-daemon warns.
respawn_period=60
respawn_max=10
pidfile="/run/arbor/web.pid"
output_log="/var/log/arbor/web.log"
error_log="/var/log/arbor/web.log"

env_file="/etc/arbor/arbor.env"

depend() {
    need arbor-daemon
    after net
}

start_pre() {
    checkpath --directory --owner root:arbor --mode 0750 /run/arbor
    checkpath --directory --owner arbor:arbor --mode 0750 /var/log/arbor
    touch /var/log/arbor/web.log
    chown arbor:arbor /var/log/arbor/web.log
    chmod 0640 /var/log/arbor/web.log
    if [ ! -x /usr/bin/setpriv ]; then
        eerror "setpriv not found (sys-apps/util-linux). Cannot apply web hardening."
        return 1
    fi
}
